# system_files/usr/lib/systemd/user/noctalia.service
[Unit]
Description=Noctalia v5 desktop shell (Mazzite)
Documentation=https://docs.noctalia.dev/v5/
# Gate the unit on the Mazzite MangoWM session. Bazzite's plasmalogin
# (renamed SDDM) takes XDG_SESSION_DESKTOP from the `DesktopNames=`
# field of the .desktop file, NOT the filename. Our
# mangowm-noctalia.desktop has `DesktopNames=MangoWM`, so the env var
# that arrives in the user session is `MangoWM`, not `mangowm-noctalia`.
# Use that value here so the unit fires when the user picks our session
# — and does NOT fire in the default Bazzite (KDE/Plasma/GNOME) session.
ConditionEnvironment=XDG_SESSION_DESKTOP=MangoWM
# Belt-and-suspenders: only run inside a Wayland session.
# (ConditionEnvironment is a [Unit] directive, NOT [Service]. systemd
# v256+ silently ignores it in [Service] and warns at parse time.)
ConditionEnvironment=WAYLAND_DISPLAY
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/noctalia
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
