# system_files/usr/lib/systemd/user/noctalia.service
[Unit]
Description=Noctalia v5 desktop shell (Mazzite)
Documentation=https://docs.noctalia.dev/v5/
# Gate the unit on the Mazzite MangoWM session. Bazzite's plasmalogin
# (renamed SDDM) takes XDG_SESSION_DESKTOP from the `DesktopNames=`
# field of the .desktop file, NOT the filename. Our
# mangowm-noctalia.desktop has `DesktopNames=MangoWM`, so the env var
# that arrives in the user session is `MangoWM`, not `mangowm-noctalia`.
# Use that value here so the unit fires when the user picks our session
# — and does NOT fire in the default Bazzite (KDE/Plasma/GNOME) session.
#
# We do NOT additionally gate on WAYLAND_DISPLAY here. Plasmalogin
# doesn't reliably propagate WAYLAND_DISPLAY into the user manager at
# the moment of evaluation, so the gate fails intermittently with
# "unmet condition check ConditionEnvironment=WAYLAND_DISPLAY" even
# though the var is present shortly after. XDG_SESSION_DESKTOP=MangoWM
# is the real discriminator, and `WantedBy=graphical-session.target` +
# mazzite-session.service activating the target handle the timing.
ConditionEnvironment=XDG_SESSION_DESKTOP=MangoWM
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/noctalia
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
